# bootstrap/swarm-bootstrap.service
[Unit]
Description=SwarmDrop Bootstrap + Relay Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/swarm-bootstrap
ExecStart=/opt/swarm-bootstrap/swarm-bootstrap run \
    --tcp-port 4001 \
    --quic-port 4001
Restart=on-failure
RestartSec=5

# 日志
Environment=RUST_LOG=info
StandardOutput=journal
StandardError=journal
SyslogIdentifier=swarm-bootstrap

# 资源限制
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
